# gforge-plugin-scmcvs/rpm-specific/cron.d/gforge-plugin-scmcvs
#
# GForge
# Copyright 2005 (c) Guillaume Smet
# http://gforge.org/
#
# @version   $Id: gforge,v 1.2 2004/12/07 16:12:24 gsmet Exp $
#

# Change path to the GForge directory.
GFORGE="/usr/share/gforge"
GFORGE_PLUGINS="/usr/lib/gforge/plugins"

# Change paths to the GForge directory.
# Don't reference $GFORGE because this won't work!
# You may need to change the pathname to php CLI (command line interface)
PHP="/usr/bin/php -q -d include_path=.:/etc/gforge:/usr/share/gforge:/usr/share/gforge/www/include"

# Don't mail command output to anyone.
MAILTO=""

#
#       Generate CVS tarballs
#
5 2 * * * root $PHP $GFORGE_PLUGINS/scmcvs/cronjobs/tarballs.php