# src/etc/cron.d/fusionforge-web-vhosts
37 7,19 * * * root forge_run_job create-vhosts.php
