# src/etc/cron.d/fusionforge-web-vhosts
PATH=@bindir@:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
37 7,19 * * * root forge_run_job web-vhosts/create-vhosts.php
