# src/etc/cron.d/fusionforge-shell
#
# Regular cron jobs for shell-postgresql
#
PATH=@bindir@:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# SCM/user/group update
15 * * * * root forge_run_job shell/homedirs.php

# SSH keys
20 * * * * root forge_run_job shell/ssh_create.php
