# src/plugins/globalsearch/packaging/cron.d/plugin-globalsearch
30 3 * * * gforge [ -x @PLUGINS_PATH@/globalsearch/bin/get_assoc_sites_projects.pl ] && @PLUGINS_PATH@/globalsearch/bin/get_assoc_sites_projects.pl > /dev/null 2>&1
