# src/etc/init.d/fusionforge-db-remote
#!/bin/bash
### BEGIN INIT INFO
# Provides:          postgresql
# Required-Start:    $local_fs $remote_fs $network $syslog $named
# Required-Stop:     $local_fs $remote_fs $network $syslog $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: Fake init dependency for remote DB installs
### END INIT INFO

true

# Rationale:
# $ update-rc.d fusionforge-systasksd defaults
# insserv: Service postgresql has to be enabled to start service fusionforge-systasksd
# insserv: exiting now!
# update-rc.d: error: insserv rejected the script header
