# src/plugins/scmhook/etc/cron.d/fusionforge-plugin-scmhook
# Apply hooks every hour
0 * * * * root forge_run_plugin_job scmhook updateScmRepo.php
