# gforge/packaging/cron.d/backup1
# Daily site backup
0 2 * * *  root $PHP $FFCRON/backup_site.php

