# gforge/packaging/cron.d/plugin-mediawiki
#
# Regular cron jobs for plugin-mediawiki
#

# Create new wikis four times an hour
*/15 * * * * root [ -x $FFDIR/plugins/mediawiki/cronjobs/create-wikis.php ] && $FFDIR/plugins/mediawiki/cronjobs/create-wikis.php
