# gforge/plugins/eirc/debian/cron.d
0 0 * * * gforge [ -x /usr/lib/sourceforge/plugins/eirc/bin/eirc ] && /usr/lib/sourceforge/plugins/eirc/bin/eirc > /dev/null 2>&1
