# gforge-plugin-helloworld/debian/cron.d
0 0 * * * sourceforge [ -x /usr/lib/sourceforge/plugins/helloworld/bin/sayhello ] && /usr/lib/sourceforge/plugins/helloworld/bin/sayhello > /dev/null 2>&1
