# plugins/cvstracker/rpm-specific/cron.d/gforge-plugin-cvstracker
#
# GForge
# Copyright 2005 (c) Guillaume Smet
# http://fusionforge.org/
#

# Change path to the GForge directory.
GFORGE="/usr/share/gforge"
GFORGE_PLUGINS="/usr/share/gforge/plugins"

# Change paths to the GForge directory.
# Don't reference $GFORGE because this won't work!
# You may need to change the pathname to php CLI (command line interface)
PHP="/usr/bin/php -q -d include_path=.:/etc/gforge:/usr/share/gforge:/usr/share/gforge/www/include"

# Don't mail command output to anyone.
MAILTO=""

#
#       Updates the loginfo files to include CVS Tracker
#
30 * * * * root $PHP $GFORGE_PLUGINS/cvstracker/bin/update_loginfo.php
