# gforge/plugins/mediawiki/etc/cron.d/mediawiki
# Create new mediawikis four times an hour
*/15 * * * * $FFUSER $PHP $FFDIR/plugins/mediawiki/cronjobs/create-wikis.php
