# Branch_4_5/gforge/plugins/cvstracker/debian/cron.d
35 1 * * * root [ -x /usr/lib/gforge/plugins/cvstracker/bin/update_loginfo.php ] && /usr/lib/gforge/plugins/cvstracker/bin/update_loginfo.php -d include_path=/etc/gforge:/usr/share/gforge/:/usr/share/gforge/www/include
