# src/plugins/repositoryapi/etc/cron.d/fusionforge-plugin-repositoryapi
# Gather activities from SCM repositories
0 */3 * * * root forge_run_plugin_job repositoryapi parse_scm_repo_activities.php
