# src/plugins/repositoryapi/etc/cron.d/fusionforge-plugin-repositoryapi
PATH=@bindir@:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Gather activities from SCM repositories
0 */3 * * * root forge_run_plugin_job repositoryapi parse_scm_repo_activities.php
