# src/packaging/cron.d/backup1
# Daily site backup
0 2 * * *  root forge_run_job backup_site.php

