# src/packaging/cron.d/dns-bind9
#
# Regular cron jobs for the dns-bind9 package
#

# DNS Update
0 * * * * root [ -f /var/lib/@OLDPACKAGE@/bind/dns.head ] && [ -x @BINARY_PATH@/dns_conf.pl ] && @BINARY_PATH@/dns_conf.pl &&  /usr/sbin/invoke-rc.d --quiet bind9 reload 
