# src/packaging/cron.d/homedirs
#
# IMPORTANT!! READ THE README.root file before continuing!
#

# Create the shell users accounts and groups home directories
15 * * * * root forge_run_job homedirs.php

