# src/packaging/cron.d/lists-mailman
#
# Regular cron jobs for the lists-mailman package
#

# Mailing-list creation
55 * * * * root [ -x /usr/share/gforge/bin/create-mailing-lists.pl ] && /usr/share/gforge/bin/create-mailing-lists.pl > /dev/null 2>&1
