# src/packaging/cron.d/lists-mailman
#
# Regular cron jobs for the lists-mailman package
#

# Mailing-list creation
55 * * * * root $PHP $FFCRON/mail/mailing_lists_create.php
