# src/packaging/cron.d/mail
# Create the new mailing lists
06 * * * * root forge_run_job mail/mailing_lists_create.php

# Create the /etc/aliases file
08 * * * * root forge_run_job mail/mailaliases.php

# Send email to users with pending tasks
#30 5 * * 1,2,3,4,5 root forge_run_job daily_task_email.php

