# src/packaging/cron.d/shell-postgresql
#
# Regular cron jobs for shell-postgresql
#

# SCM/user/group update
0 * * * * root [ -x @BINARY_PATH@/update-user-group-ssh.sh ] && @BINARY_PATH@/update-user-group-ssh.sh > /dev/null 2>&1
