# src/packaging/cron.d/web-apache2-vhosts
#
# Regular cron jobs for the web-apache2-vhosts package
#

# Enable the virtual hosts
37 7,19 * * * root [ -x @BINARY_PATH@/create-vhosts.sh ] && @BINARY_PATH@/create-vhosts.sh > /dev/null 2>&1
