# src/plugins/quotamanagement/etc/cron.d/fusionforge-plugin-quotamanagement
#
# Regular cron jobs for the plugin-quotamanagement package
#
PATH=@bindir@:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# calculate quota every night
55 2 * * * root forge_run_plugin_job quotamanagement quotamanagement.php 
